# audio-led.service
[Unit]
Description=Audio LED Controller
After=pipewire.service wireplumber.service sound.target
Requires=wireplumber.service

[Service]
Type=simple
ExecStart=/usr/local/bin/audio-led.sh
Restart=on-failure
RestartSec=5

# Environment variables for user PulseAudio (PipeWire) session
Environment="PULSE_RUNTIME_PATH=/run/user/%U/pulse/"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus"

# Delay to ensure audio stack is fully ready
ExecStartPre=/bin/sleep 3

[Install]
WantedBy=default.target
